# infra/deploy-app1-autopull.timer
[Unit]
Description=Poll GitHub and deploy jamafurn.com every minute

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
