# infra/deploy-yiwu-wax-autopull.service
[Unit]
Description=Pull main and deploy yiwuwax.com (yiwu-wax)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/deploy-yiwu-wax-autopull.sh
TimeoutStartSec=2400
Nice=10
